# travel-points.timer
[Unit]
Description=Travel Points daily digest (00:00 UTC)

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
